# processPaymentOnLive.ini
[program:processPaymentOnLive]
command=/usr/bin/python3 /var/www/vhosts/rmuictonline.com/ussd.rmuictonline.com/processPaymentOnLive.py
autostart=true
autorestart=true
stderr_logfile=/var/www/vhosts/rmuictonline.com/ussd.rmuictonline.com/processPaymentOnLiveErr.log
stdout_logfile=/var/www/vhosts/rmuictonline.com/ussd.rmuictonline.com/processPaymentOnLiveOut.log
